# AeraSite/pwi/modules/top/etc/crontab
0 1 * * * /PWServer/scripts/pkupdate
5 1 * * * /PWServer/scripts/pkupdate2
10 1 * * * /PWServer/scripts/pkupdate3